# installer/solar-controller.service
[Unit]
Description=Solar-controller
After=network.target auditd.service

[Service]
ExecStart=solar-controller
Restart=always
